# udev/rules.d/99-odroid-wiringpi-pwm.rules
ACTION=="add", KERNEL=="pwmchip*", SUBSYSTEM=="pwm", RUN+="/bin/sh /etc/odroid-allowuser.sh /sys/class/pwm/pwmchip*"
